# BME280-InfluxDB.timer
[Unit]
Description=BME280-InfluxDB.js

[Timer]
OnBootSec=3min
OnUnitActiveSec=20s

[Install]
WantedBy=timers.target
